# data/docker/supervisord.conf
[supervisord]
nodaemon=false

[program:sdServer]
command=conda run -n ldm python '/modules/stable-ui/sd_model_server.py'
autostart=true
autorestart=false
stdout_logfile=/logs/docker_server.out
stderr_logfile=/logs/docker_server.err
