# data/docker/supervisord.conf
[supervisord]
nodaemon=false

[program:sdServer]
command=/home/sd-server.sh
autostart=true
autorestart=false
stdout_logfile=/logs/docker_server.out
stderr_logfile=/logs/docker_server.err
